# example-units/example.service
[Unit]
Description=Tweet a new picture
After=network.target

[Service]
ExecStart=flickr_image_bot
